# docker/s6/nginx/run
#!/usr/bin/env sh
exec /usr/sbin/nginx -g "daemon off; error_log /dev/stderr error;"