# docker/s6/php-fpm/run
#!/usr/bin/env sh
set -o pipefail
exec /usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.d/www.conf